# init_scripts/jetbot_stats/jetbot_telemetria.timer
[Unit]
Description=Timer for jetbot_telemetria

[Timer]
OnBootSec=30sec
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
